# ops/launchd/com.clawberto.aerodrome.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clawberto.aerodrome.watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/marko/.openclaw/workspace/Clawberto-Aerodrome/scripts/aerodrome-lp-watcher-launchd.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/marko/.openclaw/workspace/Clawberto-Aerodrome</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>15</integer>
  <key>StandardOutPath</key>
  <string>/Users/marko/.openclaw/workspace/Clawberto-Aerodrome/runs/aerodrome-lp-supervisor/watcher.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/marko/.openclaw/workspace/Clawberto-Aerodrome/runs/aerodrome-lp-supervisor/watcher.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/marko/.foundry/bin</string>
    <key>HERMES_RPC_URL</key>
    <string>https://base-rpc.publicnode.com</string>
    <key>HERMES_LP_EXECUTE</key>
    <string>0</string>
    <key>HERMES_DISCOVERY_FORWARD_SCAN</key>
    <string>0</string>
  </dict>
</dict>
</plist>
